# docker/supervisord.conf
; Single-container process tree for the aMaze platform.
;
; Process order (priority lower = starts earlier):
;   10  redis                  — state store; other services fail-closed on
;                                Redis outage (per CLAUDE.md §8).
;   20  orchestrator           — FastAPI: /register, /resolve/*, /health.
;                                Must be up before the proxy so the policy
;                                addon can bootstrap (load YAMLs, open
;                                Redis connection).
;   30  proxy                  — mitmproxy + policy addon. Last so that
;                                downstream dependencies are already warm.
;
; `autorestart=true` everywhere: a crashed process should come back. Logs
; go to /var/log/amaze/{name}.log so `docker logs` + host bind-mount give
; operators a unified view.

[supervisord]
nodaemon=true
logfile=/var/log/amaze/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; --- redis ----------------------------------------------------------------
[program:redis]
command=/usr/bin/redis-server --appendonly yes --dir /data/redis --bind 127.0.0.1 --protected-mode no
priority=10
autorestart=true
stdout_logfile=/var/log/amaze/redis.log
stderr_logfile=/var/log/amaze/redis.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; --- orchestrator ---------------------------------------------------------
[program:orchestrator]
command=/usr/local/bin/uvicorn services.orchestrator.main:app --host 0.0.0.0 --port 8001 --no-access-log
directory=/app
environment=PYTHONPATH="/app",REDIS_URL="redis://127.0.0.1:6379",CONFIG_DIR="/app/config"
priority=20
autorestart=true
stdout_logfile=/var/log/amaze/orchestrator.log
stderr_logfile=/var/log/amaze/orchestrator.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; --- proxy (mitmproxy + addons) ------------------------------------------
[program:proxy]
command=/usr/local/bin/mitmdump
    --listen-host 0.0.0.0
    --listen-port 8080
    --set confdir=/opt/mitmproxy
    --set block_global=false
    -s services/proxy/main.py
directory=/app
environment=PYTHONPATH="/app",REDIS_URL="redis://127.0.0.1:6379",ORCHESTRATOR_URL="http://127.0.0.1:8001",CONFIG_DIR="/app/config"
priority=30
autorestart=true
stdout_logfile=/var/log/amaze/proxy.log
stderr_logfile=/var/log/amaze/proxy.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
